# base/init.d/remount
#!/sbin/openrc-run
command=/usr/bin/remount